# rtk_openwrt_src/build_dir/target-mips-rlx4181-linux/root-realtek/usr/lib/opkg/info/rtk_app.list
/usr/sbin/wpsled
/usr/sbin/rtk_txcalr
/usr/sbin/flash
/usr/sbin/rtk_pcap
/etc/init.d/rtk_app
